# cricket-webhook.service
[Unit]
Description=Cricket Live Webhook Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/live-ipl-updates
Environment="WEBHOOK_URL=https://poke.com/api/v1/inbound/webhook"
Environment="WEBHOOK_AUTH=Bearer eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJzdWIiOiI4NjBmOWZjZS1jZTNhLTQ0NzEtOTUzOS02ZTZjOGRhODEwNGIiLCJqdGkiOiI0ZDY2MTBjMy03ZTczLTQzMGItODc1Yy1jNzUxMmUxMWI4NmYiLCJpYXQiOjE3NzQ4ODQ2MzgsImV4cCI6MjA5MDI0NDYzOH0.Z0OWRx1jThCNNCloF1QiB9-a8hHtSQNLHA7tzIYGubs"
Environment="PORT=6666"
ExecStart=/usr/bin/python3 proxy-python/server.py
Restart=always
RestartSec=10
StandardOutput=append:/tmp/cricket_server.log
StandardError=append:/tmp/cricket_server.log

[Install]
WantedBy=multi-user.target
